# .config/systemd/user/abra-session.target
[Unit]
Description = ABRA – combined graphical session target
Requires    = picom.service polybar.service autorandr.service
